# packages/rpi2casterd/rpi2casterd-2.3.2.tar.gz/rpi2casterd-2.3.2/data/rpi2casterd.service
[Unit]
Description=rpi2casterd - hardware interface control daemon for rpi2caster
After=network.target

[Service]
ExecStart=/usr/local/bin/rpi2casterd
User=monotype
Group=monotype
Restart=on-abort

[Install]
WantedBy=multi-user.target
